# lib/capistrano/templates/systemd/master.target.erb
# This file is generated by capistrano-procfile (at <%= Time.now.utc %>).
#
# Do not manually change this file.
[Unit]
Description=<%= @application_name %> services target (generated by capistrano-procfile)
Requires=<%= @services_filenames.join(" ") %>
After=syslog.target network.target

[Install]
WantedBy=multi-user.target
