# lib/capistrano/templates/systemd/master.target.erb
# This file is generated by capistrano-procfile (at <%= Time.now.utc %>).
#
# Do not manually change this file.
[Unit]
Description=<%= app %> services target (generated by capistrano-procfile)
# Wants=<%= processes_names.join(" ") %>
Requires=<%= processes_names.join(" ") %>
After=syslog.target network.target

[Install]
WantedBy=multi-user.target
