# lib/capistrano/templates/systemd/process.service.erb
# This file is generated by capistrano-procfile (at <%= Time.now.utc %>).
#
# Do not manually change this file.
[Unit]
Description=<%= procname %> service of <%= app %> application (generated by capistrano-procfile)
PartOf=<%= app %>.target

[Service]
User=<%= user %>
<% unless group.nil? -%>
Group=<%= group %>
<% end -%>
<% unless umask.nil? -%>
UMask=<%= umask %>
<% end -%>

WorkingDirectory=<%= root %>
<% env_vars.each_pair do |var, env| -%>
Environment="<%= var.upcase %>=<%= env %>"
<% end -%>
ExecStart=/bin/bash -lc 'exec <%= command %>'
# ExecStop=-/sbin/start-stop-daemon --quiet --stop --retry TERM/5 --pid $MAINPID # From `nginx.service`
RestartSec=<%= restart_sec %>
Restart=<%= restart_method %>

StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=<%= syslog_identifier %>

KillMode=mixed
# If ExecStop is defined.
# KillMode=control-group
TimeoutStopSec=<%= timeout %>
