# lib/capistrano/templates/systemd/process.service.erb
# This file is generated by capistrano-procfile (at <%= Time.now.utc %>).
#
# Do not manually change this file.
[Unit]
Description=<%= @service_name %> service of <%= @application_name %> application (generated by capistrano-procfile)
PartOf=<%= @target_name %>.target

[Service]
User=<%= @user %>
<% unless @group.nil? -%>
Group=<%= @group %>
<% end -%>
<% unless @umask.nil? -%>
UMask=<%= @umask %>
<% end -%>

WorkingDirectory=<%= @root_path %>
<% @env_vars.each_pair do |var, env| -%>
Environment="<%= var.upcase %>=<%= env %>"
<% end -%>
ExecStart=/bin/bash -lc 'exec <%= @command %>'
RestartSec=<%= @restart_sec.to_s %>
Restart=<%= @restart_method.to_s %>

StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n

KillMode=mixed
TimeoutStopSec=<%= @timeout %>
